# freebsd/unit
#!/bin/sh

# PROVIDE: [[[MINECRAFT_SERVER_SERVICE_NAME]]]
# REQUIRE: LOGIN FILESYSTEMS
# KEYWORD: shutdown

. /etc/rc.subr
. [[[SERVICE_CONFIG_DIR]]]/[[[MINECRAFT_SERVER_SERVICE_NAME]]]
. [[[SERVICE_LIB_DIR]]]/FUNC_STOP

name="[[[MINECRAFT_SERVER_SERVICE_NAME]]]"
rcvar="[[[MINECRAFT_SERVER_SERVICE_NAME]]]_enable"
[[[MINECRAFT_SERVER_SERVICE_NAME]]]_user=${MINECRAFT_SERVER_EXECUTE_USER}
[[[MINECRAFT_SERVER_SERVICE_NAME]]]_group=${MINECRAFT_SERVER_EXECUTE_GROUP}
[[[MINECRAFT_SERVER_SERVICE_NAME]]]_chdir=${MINECRAFT_SERVER_ROOT}

command=[[[SCREEN_PATH]]]
command_args="-AmdS ${SCREEN_NAME} [[[JVM_PATH]]] -server -Xms${MIN_MEMORY} -Xmx${MAX_MEMORY} -jar ${JAR_PATH} nogui"

stop_cmd="FUNC_STOP"

load_rc_config $name

run_rc_command "$1"
